# target/linux/lantiq/xrx330/P02003/base-files/etc/init.d/enablePHY
#!/bin/sh /etc/rc.common

START=62

start() {
	echo "enable PHY..."
	switch_cli IFX_ETHSW_PORT_LINK_CFG_SET nPortId=1 bLinkForce=1 eLink=0
	switch_cli IFX_ETHSW_PORT_LINK_CFG_SET nPortId=2 bLinkForce=1 eLink=0
	switch_cli IFX_ETHSW_PORT_LINK_CFG_SET nPortId=3 bLinkForce=1 eLink=0
	switch_cli IFX_ETHSW_PORT_LINK_CFG_SET nPortId=4 bLinkForce=1 eLink=0
}
