# target/linux/lantiq/xrx330/P02003/base-files/etc/init.d/lanStart
#!/bin/sh /etc/rc.common

START=20

start() {
	/etc/init.d/iptables start 
	/etc/init.d/ip6tables start 
	/etc/init.d/lan_services_bringup.sh start &
}
